# system_files/etc/systemd/system/var-mnt-containers_data.automount
[Unit]
Description=SMB Share Automount

[Automount]
Where=/var/mnt/containers_data
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target